# truckcontrol/long_ctl/src/vehcomm/dsrc/executable/ASD/wmefwdMon
#!/bin/sh /etc/rc.common
START=99
start() {
  killall -2 obuMon.sh
  /usr/local/bin/path/obuMon.sh &
}

stop() {
  killall -2 obuMon.sh
}
