# truckcontrol/system/silver/rc.local
/usr/sbin/sshd

/home/truckcontrol/test/truck_startup.sh -V Silvr
